# supervisord.conf

[supervisord]
nodaemon=true

[program:docker]
priority=10
command=dind wrapdocker
startsecs=0
exitcodes=0,1


[program:chown]
priority=20
command=/usr/local/bin/jenkins.sh

